# installers/blur-monitor_ubuntu-18.04_amd64/DEBIAN/prerm
#!/bin/bash

systemctl stop blur-monitor
systemctl disable blur-monitor

rm /etc/systemd/system/blur-monitor.service

rm -r /opt/blur-monitor
rm /usr/bin/blur-monitor
