# installers/blur-monitor_ubuntu-18.04_amd64/DEBIAN/prerm
#!/bin/bash

systemctl stop blur-monitor
systemctl disable blur-monitor
